# docker/supervisor/nodejs.ini
[program:pm2]
process_name=%(program_name)s
directory=/usr/src/app/
command=pm2 start /usr/src/app/process.json
autostart=true
autorestart=false
user=root
stdout_logfile=/var/log/supervisor/stdout_pm2.log
stderr_logfile=/var/log/supervisor/stderr_pm2.log